# pkgs/lighttpd/service/lighttpd/run
#!/bin/sh

exec /bin/lighttpd -D -f /etc/lighttpd.conf 2>&1
